# salt/srv/base/btrfs/btrfs-subvolume.mount.jinja
[Unit]
Before=docker.service
[Install]
RequiredBy=docker.service
[Mount]
What={{ device }}
Where={{ mount_point }}
Type=btrfs
Options=discard,ssd
